# config/supervisor-salt.conf
[program:salt-master]
command=/usr/bin/salt-master